# server/worker/mondaemon.conf
start on runlevel [2345]

exec /usr/local/sbin/mondaemon.sh
